# source/dist/util/rctemplates/runit/sgemaster/run
#!/usr/bin/sh
/opt/sge/default/common/sgemaster
exec pause
