# tests/5-keystrokes-to-setup-a-cluster-without-root-privilege/lib/overlay/etc/systemd/system/sge-auto.service
[Unit]
Description=Auto Installation
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/root/sge/auto.sh

[Install]
WantedBy=multi-user.target
