# source/dist/util/rctemplates/upstart/sgeexecd.conf
description "Grid Engine execd"

start on (remote-filesystems and runlevel [2345])
stop on runlevel [016]

respawn

pre-start script
  [ -f /etc/default/sgeexecd ] && . /etc/default/sgeexecd
end script

kill timeout 60

exec /opt/sge/$($SGE_ROOT/util/arch)/sge_execd
