# source/dist/util/rctemplates/upstart/sgemaster.conf
description "Grid Engine qmaster"

start on (net-device-up IFACE!=lo and runlevel [2345])
stop on runlevel [016]

respawn

pre-start script
  [ -f /etc/default/sgemaster ] && . /etc/default/sgemaster
end script

kill timeout 60

exec /opt/sge/$($SGE_ROOT/util/arch)
