# source/dist/util/rctemplates/runit/sgeexecd/run
#!/usr/bin/sh
/opt/sge/default/common/sgeexecd
exec pause
