# source/dist/util/rctemplates/systemd/sgeexecd.service
[Unit]
Description=Grid Engine execd
After=network-online.target remote-fs.target

[Service]
Type=forking
ExecStart=/opt/sge/default/common/sgeexecd
ExecStop=/opt/sge/default/common/sgeexecd stop
#Restart=on-failure
#RestartSec=30s

[Install]
WantedBy=multi-user.target
