# source/dist/util/rctemplates/systemd/sgemaster.service
[Unit]
Description=Grid Engine qmaster
After=network-online.target remote-fs.target autofs.service

[Service]
Type=forking
# Modify if the cell isn't "default"
#PIDFile=/opt/sge/default/spool/qmaster/qmaster.pid
ExecStart=/opt/sge/default/common/sgemaster
ExecStop=/opt/sge/default/common/sgemaster stop
#Restart=on-failure
#RestartSec=30s

[Install]
WantedBy=multi-user.target
